# deploy/systemd/memory-api.service
[Unit]
Description=Memory Empire API
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=oneshot
RemainAfterExit=true
WorkingDirectory=/home/bigphoot/memory-stack
ExecStart=/usr/bin/docker compose up -d --build memory-api
ExecStop=/usr/bin/docker compose stop memory-api
TimeoutStartSec=0

[Install]
WantedBy=default.target
